# crontab
@reboot /root/services.alemayhu.de/boot
